# cloud_monitor.conf
start on filesystem and started docker
stop on runlevel [!2345]

console log

respawn
chdir /opt/aws/aws_templates

script
	exec bash ./cloud_monitor.sh -v
end script

